# piguard/systemd/piguard.service
[Unit]
Description=PiGuard-Hybrid Self-Healing Watchdog
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /opt/piguard/daemon/piguard.py
Restart=always
RestartSec=10
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target